# com.lucas.moxie-buscador.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.lucas.moxie-buscador</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Library/Frameworks/Python.framework/Versions/3.14/bin/python3</string>
    <string>/Users/Lucas/Code/freelas/moxie_plataforma_conteudo/agente/buscar.py</string>
  </array>
  <key>StartInterval</key><integer>86400</integer>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>/Users/Lucas/Library/Logs/moxie-buscador.log</string>
  <key>StandardErrorPath</key><string>/Users/Lucas/Library/Logs/moxie-buscador.log</string>
</dict>
</plist>
